# quiz-bot-web.service
[Unit]
Description=Morning Quiz Bot Web Interface
After=network.target
# Не зависит от quiz-bot, работает независимо

[Service]
Type=simple
User=lizard
Group=lizard
WorkingDirectory=/home/lizard/morning-quiz-bot
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/lizard/morning-quiz-bot/venv/bin

ExecStart=/home/lizard/morning-quiz-bot/venv/bin/python web/run_web.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
